# crontab
*/5 * * * * $HOME/.crons/git_maintenance
*   * * * * $HOME/bin/kerb cron
